# deploy/systemd/frank-bot.service
[Unit]
Description=Frank Discord Bot
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=frankbot
Group=frankbot
WorkingDirectory=/opt/frank-bot
ExecStart=/opt/frank-bot/venv/bin/python src/bot.py
Restart=always
RestartSec=10
StartLimitInterval=350
StartLimitBurst=10

# Environment
Environment=PYTHONPATH=/opt/frank-bot
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/frank-bot/data
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Resource limits
MemoryMax=800M
CPUQuota=80%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=frank-bot

[Install]
WantedBy=multi-user.target